# deploy/arc-bots.service
[Unit]
Description=Arc Treasury Bots (Telegram + Points + Rebalance)
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/arc-treasury/bots
ExecStart=/usr/bin/tsx telegramBot.ts
Restart=on-failure
RestartSec=10
Environment=NODE_ENV=production
EnvironmentFile=/opt/arc-treasury/.env

[Install]
WantedBy=multi-user.target
